# src/resources/discord_fm.service
[Unit]
Description=Show what you're scrobbling on Discord

[Service]
ExecStart=#EXECUTABLE_PATH#

[Install]
WantedBy=default.target